# src/services/project.service
[Unit]
Description=project

[Service]
ExecStart=/path/to/project/start.sh

[Install]
WantedBy=multi-user.target
